# sh/profile.d/editor.sh
# Default editor
EDITOR="nano"
VISUAL=${EDITOR}
export EDITOR VISUAL
